# home/user/data/dotfiles/home/user/.config/systemd/user/trash-cleanup.service
[Unit]
Description=Removing trashed files
Wants=trash-cleanup.timer

[Service]
Type=simple
ExecStart=/usr/bin/trash-empty

[Install]
WantedBy=multi-user.target 
