# home/user/data/dotfiles/home/user/.config/systemd/user/trash-cleanup.timer
[Unit]
Description=Cleanup trash can

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
