# etc/udev/rules.d/69-hdparm.rules
#ACTION=="add", SUBSYSTEM=="block", KERNEL=="sda", RUN+="/usr/bin/hdparm -q -S 200 /dev/sda"
